# build/etc/dankdomain.service
[Unit]
Description=Dank Domain Door
After=local-fs.target network-online.target

[Service]
Type=simple
User=nobody
WorkingDirectory=/usr/local/games/dankdomain
Environment=TERM=linux
ExecStart=/opt/node/bin/node ./node_modules/tty.js/bin/tty.js --config ./tty/files/xterm.json
KillMode=process
TasksMax=100

[Install]
WantedBy=multi-user.target
