# build/etc/dankdomain.service
[Unit]
Description=Dank Domain Door
After=network.target httpd.service xinetd.service

[Service]
Type=simple
ExecStart=/opt/node/bin/node ./node_modules/tty.js/bin/tty.js --config ./tty/files/tty.json
KillMode=process
User=nobody
WorkingDirectory=/usr/local/games/dankdomain

[Install]
WantedBy=multi-user.target
